# dashboard-aw-main/ai-dashboard-production.service
[Unit]
Description=AI Dashboard Frontend (Production)
Documentation=https://nextjs.org/
After=network.target network-online.target
Wants=network-online.target
Requires=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/aw/dashboard-aw-main
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOSTNAME=0.0.0.0
ExecStart=/bin/bash -c 'cd /home/ubuntu/aw/dashboard-aw-main && source scripts/ensure-env.sh && npm start'
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ai-dashboard

# Security settings
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/home/ubuntu/aw/dashboard-aw-main
PrivateTmp=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Health check
TimeoutStartSec=60
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
